# etc/supervisord.d/etcd.ini
[program:etcd]
directory=/var/lib/etcd
command=/usr/bin/etcd
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

