# etc/supervisord.d/helixdns.ini
[program:helixdns]
command=/usr/bin/hdns --port 53 -forward=213.73.91.35:53
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

