# etc/supervisord.d/setup.ini
[program:setup]
command=/root/bin/setup.sh
priority=1
startretries=0
startsecs=1
redirect_stderr=true
stdout_syslog=true

