# etc/supervisord.d/startup.ini
[program:startup]
command=/root/bin/startup.sh
priority=1
startretries=0
startsecs=2
redirect_stderr=true
stdout_syslog=true
